# naboom-ops-v3.2/systemd/naboom-mqtt.service
[Unit]
Description=Naboom MQTT Subscriber (WS-first with TCP fallback)
After=nginx.service mosquitto.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=-/etc/default/naboom
WorkingDirectory=/var/www/naboomcommunity/naboomcommunity
ExecStart=/usr/local/bin/naboom-mqtt-subscriber.py
Restart=always
RestartSec=3
TimeoutStartSec=60
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/naboomcommunity /var/log/naboom
LimitNOFILE=65536
KillMode=control-group

[Install]
WantedBy=multi-user.target
